# setup/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

echo cfq > /sys/block/sda/queue/scheduler
setpci -s 0:1f.0 0xa4.w=0:1
rmmod applesmc
modprobe applesmc debug=1 smooth_led=1
/usr/local/bin/cpu-usage-leds &
